# crates/backend/staticflow-backend.service
[Unit]
Description=StaticFlow Backend API Service
After=network.target
Documentation=https://github.com/acking-you/static_flow

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/staticflow
EnvironmentFile=-/opt/staticflow/.env
ExecStart=/opt/staticflow/static-flow-backend
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 环境变量（优先级低于 .env 文件）
Environment="RUST_LOG=info"

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/staticflow/data /opt/staticflow/logs

# 资源限制
LimitNOFILE=65535
LimitNPROC=512

[Install]
WantedBy=multi-user.target
