# deployment-examples/systemd/llm-access-usage-worker.service.template
[Unit]
Description=StaticFlow LLM Access Usage Worker
After=network-online.target juicefs-llm-access.service juicefs-llm-access-usage.service
Wants=network-online.target
Requires=juicefs-llm-access.service juicefs-llm-access-usage.service

[Service]
Type=simple
User=ts_user
Group=ts_user
EnvironmentFile=/etc/llm-access/llm-access.env
EnvironmentFile=-/mnt/llm-access/config/neon.env
Environment=STATICFLOW_LOG_DIR=/var/log/staticflow-runtime
Environment=STATICFLOW_LOG_SERVICE=llm-access-usage-worker
Environment=STATICFLOW_LOG_STDOUT=0
Environment=LLM_ACCESS_USAGE_QUERY_BIND_ADDR=127.0.0.1:19081
Environment=LLM_ACCESS_STATE_ROOT=/mnt/llm-access-usage
Environment=LLM_ACCESS_DUCKDB_ACTIVE_DIR=/var/lib/staticflow/llm-access/analytics-active
# Compact work files are derived under ${LLM_ACCESS_DUCKDB_ACTIVE_DIR}/compacting
# so DuckDB rewrite work stays on local block storage instead of JuiceFS.
WorkingDirectory=/mnt/llm-access-usage
ExecStartPre=/usr/local/bin/staticflow-wait-llm-access-state analytics
ExecStart=/usr/local/bin/llm-access-usage-worker serve --bind 127.0.0.1:19081 --state-root /mnt/llm-access-usage --postgres-control-database-url-env LLM_ACCESS_CONTROL_DATABASE_URL --usage-journal-dir /var/lib/staticflow/llm-access/usage-journal --duckdb-active-dir /var/lib/staticflow/llm-access/analytics-active --duckdb-archive-dir /mnt/llm-access-usage/analytics/segments --duckdb-catalog-dir /mnt/llm-access-usage/analytics/catalog --duckdb-rollover-bytes 67108864 --usage-details-dir /mnt/llm-access-usage/details
Restart=always
RestartSec=3
NoNewPrivileges=true
StateDirectory=staticflow/llm-access
MemoryAccounting=yes
MemoryHigh=2200M
MemoryMax=3072M
MemorySwapMax=1024M
TasksMax=128
OOMPolicy=kill
# Do not add ProtectSystem/ReadWritePaths/PrivateTmp for this unit when
# /mnt/llm-access is a JuiceFS FUSE mount with default_permissions. Those
# systemd mount namespace features can fail before ExecStart with 226/NAMESPACE.

[Install]
WantedBy=multi-user.target
