# deployment-examples/systemd/llm-access.service.template
[Unit]
Description=StaticFlow LLM Access Service
After=network-online.target mnt-llm\x2daccess.mount
Wants=network-online.target
RequiresMountsFor=/mnt/llm-access

[Service]
Type=simple
User=ts_user
Group=ts_user
Environment=GEOIP_DB_PATH=/var/lib/staticflow/llm-access/geoip/GeoLite2-City.mmdb
EnvironmentFile=/etc/llm-access/llm-access.env
EnvironmentFile=-/mnt/llm-access/config/neon.env
Environment=STATICFLOW_LOG_DIR=/var/log/staticflow-runtime
Environment=STATICFLOW_LOG_SERVICE=llm-access
Environment=STATICFLOW_LOG_STDOUT=0
Environment=LLM_ACCESS_USAGE_JOURNAL_DIR=/var/lib/staticflow/llm-access/usage-journal
WorkingDirectory=/mnt/llm-access
ExecStartPre=/usr/local/bin/staticflow-wait-llm-access-state api
ExecStart=/usr/local/bin/llm-access serve --bind ${LLM_ACCESS_BIND_ADDR} --state-root ${LLM_ACCESS_STATE_ROOT} --postgres-control-database-url-env LLM_ACCESS_CONTROL_DATABASE_URL --usage-journal-dir ${LLM_ACCESS_USAGE_JOURNAL_DIR}
Restart=always
RestartSec=3
NoNewPrivileges=true
StateDirectory=staticflow/llm-access
MemoryAccounting=yes
MemoryHigh=3584M
MemoryMax=4096M
MemorySwapMax=1024M
TasksMax=256
OOMPolicy=kill
# Do not add ProtectSystem/ReadWritePaths/PrivateTmp for this unit when
# /mnt/llm-access is a JuiceFS FUSE mount with default_permissions. Those
# systemd mount namespace features can fail before ExecStart with 226/NAMESPACE.

[Install]
WantedBy=multi-user.target
